# scripts/systemd/new-voice-agent.service
[Unit]
Description=NEW-VOICE Voice Agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/new-voice
Environment="PATH=/root/new-voice/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/root/new-voice/venv/bin/python -m src.voice_agent.simple_agent dev
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
