# scripts/systemd/new-voice-agent.service
[Unit]
Description=NEW-VOICE Voice Agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/new-voice
Environment="PATH=/root/new-voice/venv/bin"
Environment="SCENARIO_PATH=examples/scenarios/salon_scenario.yaml"
ExecStart=/root/new-voice/venv/bin/python -m src.voice_agent.scenario_voice_agent dev
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
